# deploy/trading_bot.service
[Unit]
Description=V28 Production Trading Bot (legacy filename — use v28_trading_bot.service)
Documentation=https://github.com/taglia21/Algebraic-Topology-Neural-Net-Strategy
After=network.target network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=tradingbot
Group=tradingbot
WorkingDirectory=/opt/trading-bot
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONDONTWRITEBYTECODE=1"
Environment="V28_ENV=production"
Environment="V28_VERSION=28.0.0"
EnvironmentFile=/opt/trading-bot/.env

# Main execution command — points to the V28 runner (NOT legacy production_launcher)
ExecStart=/opt/trading-bot/.venv/bin/python run_v28_production.py --mode=paper

# Graceful shutdown handling
ExecStop=/bin/kill -SIGTERM $MAINPID
TimeoutStopSec=60
KillMode=mixed
KillSignal=SIGTERM

# Restart configuration
Restart=always
RestartSec=30
RestartPreventExitStatus=0

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/trading-bot/logs
ReadWritePaths=/opt/trading-bot/state
ReadWritePaths=/opt/trading-bot/cache
ReadWritePaths=/opt/trading-bot/results
ReadWritePaths=/opt/trading-bot/backups

# Resource limits (2GB RAM droplet)
MemoryMax=1536M
MemoryHigh=1280M
CPUQuota=90%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trading-bot

# Watchdog (optional - for advanced monitoring)
WatchdogSec=300

[Install]
WantedBy=multi-user.target
